# deploy/agent-drop-notify.service
# Outbox watcher → ntfy push. Install to
# /etc/systemd/system/agent-drop-notify.service and adjust paths.
[Unit]
Description=Agent Drop outbox notifier (ntfy push on new files)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/code/agent-drop
EnvironmentFile=/root/code/agent-drop/ntfy.env
ExecStart=/usr/local/bin/node /root/code/agent-drop/scripts/watch-outbox.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
