# deploy/claude-drop.service
# Web app. Install to /etc/systemd/system/claude-drop.service and adjust paths.
[Unit]
Description=Claude Drop web app (TanStack Start / Nitro)
After=network.target tailscaled.service
Wants=tailscaled.service

[Service]
Type=simple
WorkingDirectory=/root/code/claude-drop
EnvironmentFile=/root/code/claude-drop/ntfy.env
Environment=HOST=127.0.0.1
Environment=PORT=3010
Environment=NODE_ENV=production
ExecStart=/usr/local/bin/node /root/code/claude-drop/.output/server/index.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
